# roles/spider/templates/container_image_puller.service.j2
[Unit]
Description=Download container image used by spiders

[Service]
TimeoutStartSec=0
ExecStart=/usr/bin/podman pull {{ container_image }}

[Install]
WantedBy=multi-user.target
