# roles/spider/templates/container_image_puller.timer.j2
[Unit]
Description=Trigger the service to download the container image used by spiders

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
