# roles/spider/templates/spider.service.j2
[Unit]
Description=%i spider

[Service]
TimeoutStartSec=0
ExecStartPre=-/usr/bin/podman stop %i-spider
ExecStartPre=-/usr/bin/podman rm %i-spider
ExecStart=/usr/bin/podman run --env-file {{ diario_configdir }}/env  --name %i-spider --volume {{ host_data_path }}/%i:/mnt/data --user {{ containeruser }} {{ container_image }}  bash -c "cd data_collection && scrapy crawl %i"

[Install]
WantedBy=multi-user.target
