# roles/spider/templates/spider.timer.j2
[Unit]
Description=Run %i spider

[Timer]
OnCalendar={{timer.oncalendar}}
RandomizedDelaySec={{timer.randomizeddelaysec}}
Persistent=true

[Install]
WantedBy=timers.target
